# net-p2p/Sia/files/sia.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

name="Sia daemon"
command="/usr/bin/siad"
command_args="${SIA_OPTS:-"-d /var/lib/sia"}"
command_user="${command_user:-sia:sia}"
command_background=yes
retry="${retry:-TERM/30/KILL/5}"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	use net
}
